# deploy/systemd/mapit-web.service
[Unit]
Description=MapIT Web (Next.js)
After=network.target
Wants=network.target

[Service]
Type=simple
User=mapit
Group=mapit
WorkingDirectory=/opt/mapit
EnvironmentFile=/opt/mapit/apps/web/.env.local
Environment=PORT=3001
ExecStart=/usr/bin/npm --workspace=apps/web run start:prod
Restart=always
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=30

# Security hardening (basic)
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/mapit

[Install]
WantedBy=multi-user.target
